# submission/week-01/07.service
[Unit]
Description= blah

[Service]
Type=oneshot
ExecStart=/bin/echo "Hello World"